# script/reload_confd.crontab
*/1 * * * * cd /home/confd/ && ./confd -c config.json -s reload
